# ckan/crontabs/ckan
# It will be run under the <filename> user
* * * * * /usr/bin/paster --plugin=ckanext-harvest harvester run --config=/srv/app/production.ini
